# ansible/roles/mistral_rs/templates/mistral.service.j2
[Unit]
Description=Mistral.rs Inference Server
After=network.target

[Service]
Type=simple
User={{ service_user }}
Group={{ service_group }}
WorkingDirectory={{ service_home }}
Environment=RUST_LOG=info
ExecStart={{ service_home }}/mistral-rs serve \
    --model {{ mistral_model_path }} \
    --host {{ service_host }} \
    --port {{ service_port }}
Restart=always
RestartSec=5

# Security hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
